# chef-repo-awesome-customers/cookbooks/awesome_customers_ubuntu/templates/supervisord.conf.erb
[supervisord]
nodaemon=true

[program:apach2]
command=/etc/init.d/apache2-customers start
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mysql]
command=/etc/init.d/mysql-default start
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
